# deploy/systemd/speak-up-frontend.service
[Unit]
Description=Speak Up Next.js frontend
After=network.target speak-up-backend.service

[Service]
Type=simple
WorkingDirectory=/opt/speak_up/frontend
Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_API_BASE_URL=
ExecStart=/usr/bin/env npm run start -- --hostname 127.0.0.1 --port 3000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
