# deploy/systemd/speak-up-backend.service
[Unit]
Description=Speak Up FastAPI backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/speak_up/backend
EnvironmentFile=/opt/speak_up/.env
ExecStart=/opt/speak_up/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --proxy-headers --no-access-log
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
